# rc.local
iptables-restore < /etc/iptables.ipv4.nat
sleep 10
openvpn --config /etc/openvpn/client.ovpn
service isc-dhcp-server start
exit 0
